# supervisor/conf.d/cron.conf
[program:cron]
command=crond -c /var/spool/cron/crontabs -f -L /var/log/cron/cron.log
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true