# supervisor/conf.d/filebeat.conf
[program:filebeat]
command=filebeat -e -c /etc/filebeat/filebeat.yml
autostart=true
autorestart=true
priority=1
stdout_events_enabled=true
stderr_events_enabled=true